# mercury-root/etc/systemd/system/b2-backup.service
[Service]
Type=simple
User=alex
ExecStart=
ExecStart=/usr/bin/rclone sync /home/alex/Documents/ b2-backup:alexalex-mercury-backup/documents --fast-list --transfers 32 --progress
