# neptune-root/etc/systemd/system/b2-backup.service
[Service]
Type=simple
User=alex
ExecStart=
ExecStart=/home/alex/.local/bin/run-backups
