# neptune-root/etc/systemd/system/reflector.service.d/flags.conf
[Service]
ExecStart=
ExecStart=/usr/bin/reflector \
# Set the output path where the mirrorlist will be saved (--save).
	--save /etc/pacman.d/mirrorlist \
# Select the transfer protocol (--protocol).
	--protocol https \
# Select the country (--country).
# Consult the list of available countries with "reflector --list-countries" and
# select the countries nearest to you or the ones that you trust. For example:
	--country "United States,Canada" \
	--score 10 \
# Sort the mirrors by synchronization time (--sort).
	--sort rate
# See "reflector --help" for details.
