# root/etc/systemd/system/b2-backup.timer
[Unit]
Description=Run b2-backup daily

[Timer]
OnCalendar=*-*-* 13:00:00
Persistent=true

[Install]
WantedBy=timers.target
